# kubetest2-ec2/config/ubuntu/containerd.target
[Unit]
Description=Containerd
[Install]
WantedBy=multi-user.target
